# systemd/flatline-short-test.timer
[Unit]
Description=Flatline SMART short self-test (weekly)

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
